# systemd/transithub.service
[Unit]
Description=TransitHub NYC subway LED sign
# Wait for the one-shot self-update to finish first, but only Wants= it (not
# Requires=) so the sign still starts if the update is skipped or fails.
After=network-online.target transithub-update.service
Wants=network-online.target transithub-update.service

[Service]
Type=simple
WorkingDirectory=/home/pi/transithub
ExecStart=/home/pi/transithub/.venv/bin/transithub --config /home/pi/transithub/config.yaml
Restart=always
RestartSec=5
# Matrix timing needs root, and config's matrix.drop_privileges: false keeps it
# so the app can still read the venv/fonts under the install dir.
User=root

[Install]
WantedBy=multi-user.target
